# debian/extra/dbus.service
[Unit]
Description=D-Bus System Message Bus (deputy systemd dummy)

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=true
