# debian/systemd.prerm
#! /bin/sh
set -e

# we must stop the deputy init on upgrades to later Ubuntu releases which have a real systemd
if [ "$1" = "upgrade" ] && dpkg --compare-versions "$2" gt 219; then
    invoke-rc.d systemd stop || true
fi

#DEBHELPER#

