# units/systemd-resolved-update-resolvconf.path
[Path]
PathChanged=/run/systemd/resolve/stub-resolv.conf
